# deploy/stable/systemd/rgbt-model.user.service.template
[Unit]
Description=RGBT UAV model service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=DEPLOY_CONFIG=__DEPLOY_CONFIG__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/scripts/start_model_stable_local.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
